# packaging/ngeo.sh
#! /bin/sh
##############################################################################
# @copyright GMV 2013. Property of GMV; all rights reserved
# @project NGEO T4
# @version $Rev:  $
# @date $LastChangedDate: 2013-09-02 14:20:00 +0200 (lun, 02 Sep 2013) $
# @purpose NGEO service script. To be placed in /etc/init.d
##############################################################################
### BEGIN INIT INFO
# Provides: ngeo service
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 5
# Default-Stop: 0 1 4 6
# Short-Description: NGEO service script
# Description: To start, stop, restart NGEO service 
### END INIT INFO

# ----------------------------------------------------------------------------
# Configuration section
# ----------------------------------------------------------------------------
# 
INTERFACE="eth0"
# Broadcast address
BCASTADDR="192.168.151.255"

# ----------------------------------------------------------------------------
# End of configuration section
# ----------------------------------------------------------------------------

# ----------------------------------------------------------------------------
# Start
# ----------------------------------------------------------------------------
ngeo_start() {  
    echo "------------------------------------------------------------------------------"
    echo " NGEO Service Start"
    echo "------------------------------------------------------------------------------"
	service httpd start
}

# ----------------------------------------------------------------------------
# Start
# ----------------------------------------------------------------------------
ngeo_stop() {
    echo "------------------------------------------------------------------------------"
    echo " NGEO Service Stop"
    echo "------------------------------------------------------------------------------"
	service httpd stop

}

# ----------------------------------------------------------------------------
# Start
# ----------------------------------------------------------------------------
ngeo_status() {
    echo -n "Checking for service ${MYSELF}:"
	service httpd status
}

# ----------------------------------------------------------------------------
# Main
# ----------------------------------------------------------------------------
case "$1" in
start)
    ngeo_start
    ngeo_status
;;
stop)
    ngeo_stop
;;
restart)
    ngeo_stop
    ngeo_start
    ngeo_status
;;
status)
    ngeo_status
;;
*)
    echo "Usage: $0 {start|stop|restart|status}"
exit 1
;;
esac

# END ########################################################################
